# laradock/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/web/bitrix/modules/main/tools/cron_events.php >> /dev/null 2>&1
15 3 * * * laradock cd /var/www/web/ssh && /usr/bin/php send_mail_after_30_day.php > send_mail_after_30_day.txt
