# laradock/workspace/crontab/laradock
* * * * * laradock /usr/bin/php /var/www/artisan schedule:run >> /dev/null 2>&1
* * * * * laradock /usr/bin/php /var/www/bitrix/modules/main/tools/cron_events.php >> /dev/null 2>&1